# systemd/openclaw-gateway.service.d/env.conf.template
[Service]
# Drop-in override: keeps secrets out of the main .service unit so the unit
# itself can be tracked in git without leaking anything.
#
# Install:
#   sudo mkdir -p /etc/systemd/system/openclaw-gateway.service.d
#   sudo cp env.conf.template /etc/systemd/system/openclaw-gateway.service.d/env.conf
#   sudo chmod 600 /etc/systemd/system/openclaw-gateway.service.d/env.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart openclaw-gateway

# If this VPS is a dedicated sandbox, allow Claude CLI to bypass the permission
# gate. Only set this if the machine has no other user data on it.
# Environment="IS_SANDBOX=1"

# Where the gateway should read API keys / secrets from. File should be chmod 600
# and contain: COMPOSIO_API_KEY=..., GEMINI_API_KEY=..., etc.
EnvironmentFile=-/root/.openclaw/.env
